# provision/roles/actual-sync/templates/kiran-actual-sync.service.j2
[Unit]
Description=Actual Budget Sync (Up Bank)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{ actual_sync_dir }}
{% for profile in actual_sync_profiles %}
ExecStart=/bin/bash -c 'source ~/.bashrc && fnm exec -- actualsync --bank up \
{% for account in profile.accounts %}
  --accounts "{{ account.actual_account_id }}={{ account.up_account_id }}" \
{% endfor %}'
Environment="UP_USER_TOKEN={{ profile.up_user_token }}"
Environment="ACTUAL_SERVER=https://actual.fewa.app"
Environment="ACTUAL_SYNC_ID={{ actual_sync_id }}"
Environment="ACTUAL_PASSWORD={{ actual_password }}"
Environment="ACTUAL_DATA={{ actual_sync_data_dir }}"
{% endfor %}
