# provision/roles/actual-sync/templates/kiran-actual-sync.timer.j2
[Unit]
Description=Actual Budget Sync Timer (Every 15 minutes)

[Timer]
OnCalendar={{ actual_sync_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
