# provision/roles/backup/tasks/main.yml
---
- name: Configure WAL-G backups
  when: backup_enabled | bool
  block:
    - name: Assert required WAL-G backup secrets are set
      ansible.builtin.assert:
        that:
          - wal_g_version | default('3.0.1') | length > 0
          - backup_wal_g_s3_bucket | default('') | length > 0
          - backup_wal_g_s3_prefix | default('') | length > 0
          - backup_wal_g_password | default('') | length > 0
          - backup_wal_g_aws_access_key_id | default('') | length > 0
          - backup_wal_g_aws_secret_access_key | default('') | length > 0
          - backup_wal_g_postgres_password | default('') | length > 0
        fail_msg: >-
          Backup is enabled but required WAL-G backup_* secrets are missing.
          Fill provision/secrets.yml and provision/group_vars/all.yml.

    - name: Create backup directories
      ansible.builtin.file:
        path: "{{ item.path }}"
        state: directory
        owner: root
        group: root
        mode: "{{ item.mode }}"
      loop:
        - { path: "{{ backup_dir }}", mode: "0750" }
        - { path: "{{ backup_dir }}/scripts", mode: "0700" }
        - { path: "{{ backup_logs_dir }}", mode: "0750" }

    - name: Detect host architecture for WAL-G binary
      ansible.builtin.command: uname -m
      register: host_arch_raw
      changed_when: false

    - name: Set WAL-G binary name based on architecture
      ansible.builtin.set_fact:
        wal_g_binary_name: >-
          {{
            'wal-g-pg-ubuntu20.04-aarch64'
              if (host_arch_raw.stdout | trim) == 'aarch64'
              else 'wal-g-pg-ubuntu-20.04-amd64'
          }}

    - name: Download WAL-G binary for host operations
      ansible.builtin.get_url:
        url: "https://github.com/wal-g/wal-g/releases/download/v{{ wal_g_version }}/{{ wal_g_binary_name }}"
        dest: /usr/local/bin/wal-g
        mode: "0755"
        force: true

    - name: Deploy WAL-G env file
      ansible.builtin.template:
        src: wal-g.env.j2
        dest: "{{ backup_dir }}/scripts/wal-g.env"
        owner: root
        group: root
        mode: "0600"
      register: wal_g_env_deployed

    - name: Deploy WAL-G scripts
      ansible.builtin.template:
        src: "{{ item.src }}"
        dest: "{{ item.dest }}"
        owner: root
        group: root
        mode: "0700"
      loop:
        - { src: wal-g-backup.sh.j2, dest: "{{ backup_dir }}/scripts/wal-g-backup.sh" }
        - { src: wal-g-prune.sh.j2, dest: "{{ backup_dir }}/scripts/wal-g-prune.sh" }
        - { src: wal-g-check.sh.j2, dest: "{{ backup_dir }}/scripts/wal-g-check.sh" }
        - { src: wal-g-restore.sh.j2, dest: "{{ backup_dir }}/scripts/wal-g-restore.sh" }
      register: wal_g_scripts_deployed

    - name: Deploy WAL-G systemd units
      ansible.builtin.template:
        src: "{{ item.src }}"
        dest: "/etc/systemd/system/{{ item.dest }}"
        owner: root
        group: root
        mode: "0644"
      loop:
        - { src: kiran-walg-backup.service.j2, dest: kiran-walg-backup.service }
        - { src: kiran-walg-backup.timer.j2, dest: kiran-walg-backup.timer }
        - { src: kiran-walg-prune.service.j2, dest: kiran-walg-prune.service }
        - { src: kiran-walg-prune.timer.j2, dest: kiran-walg-prune.timer }
        - { src: kiran-walg-check.service.j2, dest: kiran-walg-check.service }
        - { src: kiran-walg-check.timer.j2, dest: kiran-walg-check.timer }
      register: wal_g_systemd_units

    - name: Reload systemd daemon
      ansible.builtin.systemd:
        daemon_reload: true
      when: wal_g_systemd_units.changed or wal_g_env_deployed.changed or wal_g_scripts_deployed.changed

    - name: Enable and start WAL-G timers
      ansible.builtin.systemd:
        name: "{{ item }}"
        enabled: true
        state: started
      loop:
        - kiran-walg-backup.timer
        - kiran-walg-prune.timer
        - kiran-walg-check.timer
