# provision/roles/backup/templates/kiran-restic-prune.service.j2
[Unit]
Description=Kiran VM restic retention prune

[Service]
Type=oneshot
ExecStart={{ backup_dir }}/scripts/restic-prune.sh
