# provision/roles/backup/templates/kiran-restic-check.timer.j2
[Unit]
Description=Run Kiran VM restic repository check monthly

[Timer]
OnCalendar={{ backup_check_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
