# provision/roles/vaultwarden-litestream/templates/kiran-vaultwarden-litestream-check.timer.j2
[Unit]
Description=Vaultwarden Litestream health check timer

[Timer]
OnCalendar={{ vaultwarden_litestream_check_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
